# etc/systemd/user/cmc-espeak.service
[Unit]
Description=Get Ethereum price from CMC and say it with espeak

[Service]
Type=oneshot
ExecStart=~/.local/bin/cmc_espeak.sh
